# ansible/roles/xazabd-generate-miner/templates/xazabd-generate-miner.service.j2
[Unit]
After=network.target

[Service]
User={{ xazabd_user }}
Group={{ xazabd_group }}
ExecStart=/usr/local/bin/xazabd-generate-miner.sh

[Install]
WantedBy=default.target
